# instructions/fedora/files/etc/systemd/system/warpfix.service
[Unit]
Description=Fix connection for Cloudflare WARP

[Service]
ExecStart=sh /usr/local/bin/warpfix.sh

[Install]
WantedBy=multi-user.target
